# supervisor/ipmon.conf
[program:ipmon]
command=/usr/local/bin/ipmon.sh
priority=1
redirect_stderr=true
