# supervisor/telegraf.conf
[program:telegraf]
command=telegraf --config /etc/telegraf/telegraf.conf --watch-config poll
priority=1
redirect_stderr=true
